# workspace/crontab/root
* * * * * php /var/www/n20/artisan schedule:run >> /dev/null 2>&1
